# deploy/quadlet/erp4-caddy.container
[Unit]
Description=ERP4 Caddy reverse proxy
Wants=erp4-backend.service erp4-frontend.service
After=erp4-backend.service erp4-frontend.service

[Container]
Image=docker.io/library/caddy:2.9-alpine
ContainerName=erp4-caddy
Network=erp4.network
EnvironmentFile=%h/.config/containers/systemd/erp4-caddy.env
PublishPort=0.0.0.0:80:80
PublishPort=0.0.0.0:443:443
Volume=%h/.config/containers/systemd/erp4-caddy.Caddyfile:/etc/caddy/Caddyfile:ro,Z
Volume=erp4-caddy-data:/data:Z
Volume=erp4-caddy-config:/config:Z
ReadOnly=true
Tmpfs=/tmp:size=64m,mode=1777
AutoUpdate=registry

[Service]
Restart=always

[Install]
WantedBy=default.target
